# arch/cubieboard/overwrite/etc/profile.d/codec.sh
export VDPAU_DRIVER=sunxi
export LIBVA_DRIVER_NAME=vdpau
